# systemd/sysbackup-system.service
# ═══════════════════════════════════════════════════════════════
#  LinuxGuardian — Full System Backup Service
#  Triggered by linuxguardian-system.timer (weekly)
# ═══════════════════════════════════════════════════════════════

[Unit]
Description=LinuxGuardian — Full System Backup
Documentation=https://github.com/systembackup/systembackup
After=network-online.target
Wants=network-online.target
ConditionACPower=true

[Service]
Type=oneshot
ExecStart=/usr/local/bin/linuxguardian run --system

# ── Resource Limits ───────────────────────────────────────────
Nice=19
IOSchedulingClass=idle
CPUQuota=50%
IOWeight=20
MemoryMax=1G

# ── Timeout ───────────────────────────────────────────────────
# 6 hours for a full system backup
TimeoutStartSec=21600

# ── Logging ───────────────────────────────────────────────────
StandardOutput=journal
StandardError=journal
SyslogIdentifier=linuxguardian-system

# ── Environment ───────────────────────────────────────────────
Environment=SYSBACKUP_BACKUP_TYPE=system
Environment=SYSBACKUP_LIB_DIR=/usr/local/lib/linuxguardian
Environment=SYSBACKUP_CONFIG_FILE=/etc/linuxguardian/linuxguardian.conf
