# systemd/sysbackup-monitor.service
# ═══════════════════════════════════════════════════════════════
#  LinuxGuardian — Metrics Collector Service
#  Lightweight oneshot that gathers system/backup health metrics
#  Triggered by linuxguardian-monitor.timer (every 5 minutes)
# ═══════════════════════════════════════════════════════════════

[Unit]
Description=LinuxGuardian — Metrics Collector
Documentation=https://github.com/systembackup/systembackup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/linuxguardian metrics-collect

# ── Resource Limits ───────────────────────────────────────────
# Only Nice — no CPU/IO quotas needed for this fast job
Nice=19

# ── Timeout ───────────────────────────────────────────────────
# Should finish in seconds; 60s is a generous safety net
TimeoutStartSec=60

# ── Logging ───────────────────────────────────────────────────
StandardOutput=journal
StandardError=journal
SyslogIdentifier=linuxguardian-monitor

# ── Environment ───────────────────────────────────────────────
Environment=SYSBACKUP_LIB_DIR=/usr/local/lib/linuxguardian
Environment=SYSBACKUP_CONFIG_FILE=/etc/linuxguardian/linuxguardian.conf
